# etc/supervisor/conf.d/jblog.conf
[program:jblog]
environment=PATH="/home/jeffrey/.virtualenvs/jblog/bin"
command=python /home/jeffrey/workspace/jeffrey/python-demo/jblog/wsgi.py
